# templates/loki.service.j2
[Unit]
Description=Loki Log Aggregation System
After=network.target

[Service]
Type=simple
User=loki
Group=loki
# Single-binary execution
ExecStart=/usr/local/bin/loki -config.file=/etc/loki/loki-config.yaml

# Memory and Resource Management
MemoryLimit=1G
CPUQuota=50%
Restart=on-failure

[Install]
WantedBy=multi-user.target
